# brightnessd.service
[Unit]
Description=Brightness Control Fix

[Service]
ExecStart=/opt/brightnessd.py

[Install]
WantedBy=multi-user.target 